# src/rpm.systemd/preremove.sh
systemctl stop mediaserver >/dev/null
systemctl disable mediaserver >/dev/null